# BatRack/etc/systemd/system/BatRack.service
[Unit]
Description=Evaluate sensors for active bats and trigger recordings
After=network.target raspimjpeg.service raspiwebcam.service

[Service]
ExecStart=python3 -m batrack /boot/BatRack.conf
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
