# Sensorbox/etc/systemd/system/sensorproxy.service
[Unit]
Description=Sensorproxy
After=network.target

[Service]
ExecStart=/usr/local/bin/sensorproxy --verbose --config /boot/sensorproxy.yml --metering /boot/meterings.yml
Restart=always

[Install]
WantedBy=multi-user.target
