# RadioTracking/etc/systemd/system/huaweicheck.service
[Unit]
Description=Check Huawei cellular connection

[Service]
ExecStart=/usr/bin/python3 /home/pi/huaweicheck.py
Type=oneshot
