# RadioTracking/etc/systemd/system/huaweicheck.timer
[Unit]
Description=test

[Timer]
OnBootSec=60s
OnUnitActiveSec=300s

[Install]
WantedBy=timers.target
