# RadioTracking/etc/systemd/system/signal_detect_restart.service
[Unit]
Description=signal_detect restart

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart signal_detect
