# RadioTracking/etc/systemd/system/signal_detect_restart.timer
[Unit]
Description=hourly signal_detect restart

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target