# Mesh/etc/systemd/system/batmobile.service
[Unit]
Description=batmobile: batman gateway configurator

[Service]
ExecStart=/usr/local/sbin/batmobile
Type=oneshot
StandardOutput=journal
