# Mesh/etc/systemd/system/batmobile.timer
[Unit]
Description=batmobile timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1h

[Install]
WantedBy=basic.target
